# build/files/conf/etc/systemd/system/deluge-web.service
.include /usr/lib/systemd/system/deluge-web.service

[Unit]
Requires=deluged
After=deluged
