# build/files/conf/etc/systemd/system/deluged.service
.include /usr/lib/systemd/system/deluged.service

[Unit]
Requires=network
After=network
