# build/files/conf/etc/systemd/system/nmbd.service
.include /usr/lib/systemd/system/nmbd.service

[Unit]
Requires=network
After=network
