# build/files/conf/etc/systemd/system/smbd.service
.include /usr/lib/systemd/system/smbd.service

[Unit]
Requires=network
After=network
