# build/files/conf/etc/systemd/system/winbindd.service
.include /usr/lib/systemd/system/winbindd.service

[Unit]
Requires=network
After=network
